# scripts/systemd/ggnet-frontend.service
[Unit]
Description=ggNet Frontend Service
After=network.target ggnet-backend.service
Wants=network.target

[Service]
Type=simple
User=ggnet
Group=ggnet
WorkingDirectory=/opt/ggnet/frontend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"

# Main process
ExecStart=/usr/bin/npx serve -s dist -l 3000

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ggnet/frontend

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ggnet-frontend

[Install]
WantedBy=multi-user.target

