# scripts/systemd/ggnet-backend.service
[Unit]
Description=ggNet Backend Service
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=ggnet
Group=ggnet
WorkingDirectory=/opt/ggnet/backend
Environment="PATH=/opt/ggnet/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/ggnet/backend/src"
Environment="PYTHONUNBUFFERED=1"

# Main process
ExecStart=/opt/ggnet/backend/venv/bin/uvicorn main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info \
    --access-log \
    --no-server-header

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ggnet/backend /var/log/ggnet /var/lib/ggnet

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ggnet-backend

[Install]
WantedBy=multi-user.target

