# config/systemd/ggnet-frontend.service
[Unit]
Description=ggNet Frontend Preview Server (Optional)
Documentation=https://github.com/yourusername/ggnet
After=network.target ggnet-backend.service

[Service]
Type=simple
User=ggnet
Group=ggnet
WorkingDirectory=/opt/ggnet/frontend

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Start command - Serves built frontend on port 4173
# NOTE: This is optional - Nginx can serve the static files instead
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 4173

# Restart policy
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ggnet-frontend

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

################################################################################
# NOTES:
################################################################################
# This service is OPTIONAL and typically not needed.
# 
# By default, Nginx serves the built frontend static files from:
#   /opt/ggnet/frontend/dist/
#
# This service is only useful if you want to:
# - Run frontend independently without Nginx
# - Use Vite's preview server
# - Test frontend before Nginx setup
#
# To use this service:
#   sudo systemctl enable ggnet-frontend
#   sudo systemctl start ggnet-frontend
#
# Frontend will be available at: http://localhost:4173
################################################################################

